# docs/launchd/com.openclaw.pipeline.schedule.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.pipeline.schedule</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jaredsheppard/.openclaw/workspace/pipeline/run_evening_pipeline.sh</string>
        <string>120</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jaredsheppard/.openclaw/workspace/pipeline</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>22</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jaredsheppard/.openclaw/workspace/pipeline/logs/pipeline_schedule.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jaredsheppard/.openclaw/workspace/pipeline/logs/pipeline_schedule.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jaredsheppard</string>
        <key>PATH</key>
        <string>/Users/jaredsheppard/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
