# docs/launchd/com.openclaw.pipeline.catchup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.openclaw.pipeline.catchup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jaredsheppard/.openclaw/workspace/pipeline/run_evening_catchup.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jaredsheppard/.openclaw/workspace/pipeline</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/jaredsheppard/.openclaw/workspace/pipeline/logs/pipeline_catchup.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jaredsheppard/.openclaw/workspace/pipeline/logs/pipeline_catchup.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/jaredsheppard</string>
        <key>PATH</key>
        <string>/Users/jaredsheppard/anaconda3/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
